# backend/deploy/laravel-scheduler.cron.example
* * * * * cd /var/www/iqx-intel/backend && php artisan schedule:run >> /dev/null 2>&1
